# tt-device-queue.service
[Unit]
Description=Tenstorrent device queue — serializes hardware access
After=network.target

[Service]
Type=simple
ExecStart=%h/tenstorrent/tt-device-queue/.venv/bin/python3 %h/tenstorrent/tt-device-queue/server.py
Restart=on-failure
RestartSec=2
OOMPolicy=continue
# Soft default of 1024 caused EMFILE ("Too many open files") on 2026-06-10
LimitNOFILE=65536
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
# Unbuffered stdout so job lines reach the journal promptly
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
